# usbmuxd
#!/bin/sh /etc/rc.common
START=99

stop() {
    killall usbmuxd
}

start() {
    /usr/sbin/usbmuxd
}

restart() {
    killall usbmuxd
    /usr/sbin/usbmuxd
}
